# tianqing-report-web.service
[Unit]
Description=Tianqing audit report web
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/tianqing-report-app
EnvironmentFile=/data/tianqing-audit/clickhouse/.env
ExecStart=/usr/bin/python3 /opt/tianqing-report-app/tianqing_report_web.py --host 0.0.0.0 --port 18088 --report-dir /opt/tianqing-reports --app-dir /opt/tianqing-report-app --audit-policy-file /data/tianqing-audit/config/audit_policy.json --sensitive-keywords-file /data/tianqing-audit/config/sensitive_keywords.json --exclusion-file /data/tianqing-audit/config/audit_exclusions.json --log-file /data/tianqing-audit/raw-log/tianqing.log --public-base-url https://audit.daqo.com --auth-callback-url https://ai.daqo.com/audit/auth/callback --auth-cookie-domain .daqo.com --auth-proxy-base-url http://172.88.49.60:19001 --ai-model gpt-5.5 --ai-host root@172.88.49.60 --ai-container docker-api-1
Restart=always
RestartSec=3
User=root
Group=root
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/opt/tianqing-reports /opt/tianqing-report-app /data/tianqing-audit/config /data/tianqing-audit/decrypt-records
ReadOnlyPaths=/data/tianqing-audit/raw-log /data/tianqing-audit/clickhouse/.env

[Install]
WantedBy=multi-user.target
